# debian/pve-gpu-manager/DEBIAN/postinst
#!/bin/sh
# postinst: apply PVE XPU Manager patches and enable systemd service
set -e

INDEX_TPL="/usr/share/pve-manager/index.html.tpl"
HARDWARE_PM="/usr/share/perl5/PVE/API2/Hardware.pm"
SCRIPT_TAG='<script src="/pve2/js/pve-xpu-plugin.js"></script>'
USE_STMT='use PVE::API2::Hardware::XPU;'

APPLIED=0

log() {
    logger -t pve-xpu-postinst "$@" || true
    echo "pve-xpu-postinst: $*"
}

backup_if_needed() {
    local file="$1"
    local backup="${file}.pre-xpu"
    if [ ! -f "$backup" ] && [ -f "$file" ]; then
        cp "$file" "$backup"
        log "Backed up $file -> $backup"
    fi
}

patch_index_html() {
    if [ ! -f "$INDEX_TPL" ]; then
        log "WARNING: $INDEX_TPL not found, skipping"
        return 0
    fi

    if grep -qF "$SCRIPT_TAG" "$INDEX_TPL"; then
        log "index.html.tpl already patched"
        return 0
    fi

    backup_if_needed "$INDEX_TPL"

    local tmp
    tmp=$(mktemp)
    sed "s|</head>|${SCRIPT_TAG}\n</head>|" "$INDEX_TPL" > "$tmp"
    mv "$tmp" "$INDEX_TPL"
    chmod 644 "$INDEX_TPL"
    log "Patched $INDEX_TPL: inserted script tag before </head>"
    APPLIED=1
}

patch_hardware_pm() {
    if [ ! -f "$HARDWARE_PM" ]; then
        log "WARNING: $HARDWARE_PM not found, skipping"
        return 0
    fi

    if grep -qF "PVE::API2::Hardware::XPU" "$HARDWARE_PM"; then
        log "Hardware.pm already patched"
        return 0
    fi

    backup_if_needed "$HARDWARE_PM"

    local tmp
    tmp=$(mktemp)

    # Insert use statement after other use PVE::API2::Hardware:: lines,
    # insert register_method after other register_method calls,
    # and add { type => 'xpu' } to the index array
    awk '
        # Track where to insert the use statement
        /^use PVE::API2::Hardware::/ { last_use = NR }

        # Track where to insert register_method
        /^__PACKAGE__->register_method\(\{$/ { in_register = 1 }
        in_register && /^\}\);$/ { last_register = NR; in_register = 0 }

        # Buffer all lines
        { lines[NR] = $0 }

        END {
            use_inserted = 0
            register_inserted = 0
            for (i = 1; i <= NR; i++) {
                print lines[i]
                if (i == last_use && !use_inserted) {
                    print "use PVE::API2::Hardware::XPU;"
                    use_inserted = 1
                }
                if (i == last_register && !register_inserted) {
                    print ""
                    print "__PACKAGE__->register_method({"
                    print "    subclass => \"PVE::API2::Hardware::XPU\","
                    print "    path => '"'"'xpu'"'"',"
                    print "});"
                    register_inserted = 1
                }
            }
        }
    ' "$HARDWARE_PM" > "$tmp"

    # Also add { type => 'xpu' } to the index return array
    local tmp2
    tmp2=$(mktemp)
    sed "s/{ type => 'usb' },/{ type => 'usb' },\n            { type => 'xpu' },/" "$tmp" > "$tmp2"

    mv "$tmp2" "$HARDWARE_PM"
    chmod 644 "$HARDWARE_PM"
    rm -f "$tmp"
    log "Patched $HARDWARE_PM: added XPU subclass registration"
    APPLIED=1
}

enable_service() {
    if command -v systemctl >/dev/null 2>&1; then
        systemctl daemon-reload || true
        systemctl enable pve-xpu-sriov.service || log "WARNING: Failed to enable pve-xpu-sriov.service"
        systemctl start pve-xpu-sriov.service || log "WARNING: Failed to start pve-xpu-sriov.service"
        log "Enabled and started pve-xpu-sriov.service"
    fi
}

restart_pveproxy() {
    if command -v systemctl >/dev/null 2>&1; then
        systemctl restart pveproxy || log "WARNING: Failed to restart pveproxy"
        systemctl restart pvedaemon || log "WARNING: Failed to restart pvedaemon"
        log "Restarted pveproxy and pvedaemon"
    fi
}

case "$1" in
    configure)
        # Copy default VF templates to /etc/pve/local/ if not already present
        if [ -f /usr/share/pve-xpu/xpu-vf-templates.conf ] && [ ! -f /etc/pve/local/xpu-vf-templates.conf ]; then
            cp /usr/share/pve-xpu/xpu-vf-templates.conf /etc/pve/local/xpu-vf-templates.conf 2>/dev/null || \
                log "WARNING: Could not copy templates to /etc/pve/local/ (pmxcfs may not be mounted)"
        fi
        patch_index_html
        patch_hardware_pm
        enable_service
        if [ "$APPLIED" -eq 1 ]; then
            restart_pveproxy
        fi
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

# Automatically added by dh_installsystemd/13.14.1ubuntu5
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	# The following line should be removed in trixie or trixie+1
	deb-systemd-helper unmask 'pve-xpu-sriov.service' >/dev/null || true

	# was-enabled defaults to true, so new installations run enable.
	if deb-systemd-helper --quiet was-enabled 'pve-xpu-sriov.service'; then
		# Enables the unit on first installation, creates new
		# symlinks on upgrades if the unit file has changed.
		deb-systemd-helper enable 'pve-xpu-sriov.service' >/dev/null || true
	else
		# Update the statefile to add new symlinks (if any), which need to be
		# cleaned up on purge. Also remove old symlinks.
		deb-systemd-helper update-state 'pve-xpu-sriov.service' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.14.1ubuntu5
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null || true
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		deb-systemd-invoke $_dh_action 'pve-xpu-sriov.service' >/dev/null || true
	fi
fi
# End automatically added section


exit 0
