# debian/pve-gpu-manager/lib/systemd/system/pve-xpu-sriov.service
[Unit]
Description=Apply Intel GPU SR-IOV Configuration
After=systemd-modules-load.service
After=dev-dri-card0.device
Wants=dev-dri-card0.device
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/lib/pve-xpu/apply-sriov-config.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
